# wsl2-dns-search.conf
[program:wsl2-dns-search]
command=%%WSL2_DNS_SEARCH_SCRIPT%%
